# provisioning/baikonur.service.in
[Unit]
Description=Baikonur build central service

[Service]
WorkingDirectory=SITE_DIR
Type=simple
ExecStart=PACKAGE_PREFIX/baikonur --port=PORT

[Install]
WantedBy=multi-user.target
