# provisioning/baikonur.socket.in
[Unit]
Description=Socket for baikonur connections
PartOf=baikonur.service

[Socket]
ListenStream=IP:PORT
Accept=true

[Install]
WantedBy=sockets.target
