# provisioning/baikonur.service.in
[Unit]
Description=Baikonur build central service
After=network.target baikonur.socket
Requires=baikonur.socket

[Service]
WorkingDirectory=HOME_DIR
Type=Forking
ExecStart=baikonur --port=PORT

[Install]
WantedBy=multi-user.target
